# start-dashboard.conf
start on desktop-start
stop on desktop-end

exec /home/USER/dashboards/start-dashboard.sh
